# files/systemd/system/non-nixos-gpu.service
[Unit]
Description=GPU driver setup for Nix on non-NixOS Linux systems

[Install]
WantedBy=multi-user.target

[Service]
Type=oneshot
ExecStart=ln -nsf /home/thilo/.local/share/gpu-drivers /run/opengl-driver
RemainAfterExit=yes
